# docker/orca_controller/supervisord.conf
[supervisord]
nodaemon=true

[program:orca_sm]
command=/etc/init.d/orca_controller-11080 clean-restart
autorestart=false
