# docker/orca_sm/supervisord.conf
[supervisord]
nodaemon=true

;[program:mysqld]
;command=/sbin/service mysqld start
;autorestart=false

[program:orca_sm]
command=/etc/init.d/orca_sm-14080 clean-restart
autorestart=false
