# docker/orca_am_broker/supervisord.conf
[supervisord]
nodaemon=true

[program:orca_sm]
command=/etc/init.d/orca_am+broker-12080 clean-restart
autorestart=false
